# oldCode/jainode.conf
description "node.js testing jainode service"
author "david@sabaitechnology.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

script
	export HOME="/home/vaelyn"
	exec sudo -u vaelyn /usr/local/bin/node /home/vaelyn/jairo/node/jainode.js >> /var/log/jainode.log 2>&1
end script
